# rootfs/etc/cron.hourly/mcdebian
* * * * * root /usr/sbin/free_mem_caches
